# meepo-web.service
[Unit]
Description=Meepo Web (Next.js)
After=network.target

[Service]
Type=simple
User=meepo
Group=meepo
WorkingDirectory=/home/meepo/meepo-bot/apps/web
EnvironmentFile=/etc/meepo/meepo-web.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target